# external/overlay/etc/inittab
::sysinit:/bin/mount -a
::sysinit:/bin/mkdir -p /dev/pts
::sysinit:/bin/mount -t devpts devpts /dev/pts

# Start udev daemon
::sysinit:/sbin/udevd --daemon
::sysinit:/sbin/udevadm trigger

# Start your UI after udev sets up devices
#::respawn:/usr/bin/genesis-ui

# Optional shell on tty1
::askfirst:-/bin/sh

